# homework-june-15/Karen_Yesayan/rate.service
[Service]
ExecStart=/home/kyesayan/for_index.html.sh
Restart=on-failure
