# homework-june-15/Hovsep_Petrosyan/rate.service
[Service]
ExecStart=/home/ubuntu/for_index.html.sh
Restart=on-failure
